# webservice/signing/macosx/launchd/org.eclipsefoundation.webservice.signing.macosx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
        <dict>
                <key>Label</key>
                <string>org.eclipsefoundation.webservice.signing.macosx</string>

                <key>Program</key>
                <string>/Users/genie/webservice/signing/start.sh</string>

                <key>WorkingDirectory</key>
                <string>/Users/genie/webservice/signing</string>


                <key>UserName</key>
                <string>_www</string>
                <key>GroupName</key>
                <string>_www</string>
                <key>SessionCreate</key>
                <true/>


                <key>RunAtLoad</key>
                <true/>
                <key>KeepAlive</key>
                <true/>
        </dict>
</plist>